# squashfs-root/var/lib/dpkg/info/mariadb-server-5.5.conffiles
/etc/apparmor.d/usr.sbin.mysqld
/etc/logcheck/ignore.d.paranoid/mariadb-server-5_5
/etc/logcheck/ignore.d.workstation/mariadb-server-5_5
/etc/logcheck/ignore.d.server/mariadb-server-5_5
/etc/init.d/mysql
/etc/mysql/conf.d/mysqld_safe_syslog.cnf
/etc/mysql/debian-start
/etc/logrotate.d/mysql-server
